# Library/LaunchDaemons/com.amazon.ec2.macos-init.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Debug</key>
	<true/>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:</string>
	</dict>
	<key>KeepAlive</key>
	<dict>
		<key>SuccessfulExit</key>
		<false/>
	</dict>
	<key>Label</key>
	<string>com.amazon.ec2.macosinit</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/bin/ec2-macos-init</string>
		<string>run</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/var/log/ec2-macos-init.log</string>
	<key>StandardOutPath</key>
	<string>/var/log/ec2-macos-init.log</string>
	<key>UserName</key>
	<string>root</string>
</dict>
</plist>